# MiHive.cron
*/5 * * * * python /home/pi/MiHive/MiHive.py
